# 30_config/com.notebooklm-mcp.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.notebooklm-mcp.refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>__HOME__/.local/bin/uv run notebooklm-sync --all --apply --changed-only</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/.config/notebooklm-mcp/refresh.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.config/notebooklm-mcp/refresh.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
